# private_dot_config/systemd/user/signal-i3status-rs@.service
[Service]
Type=oneshot
ExecStart=/usr/bin/pkill --signal SIGRTMIN+%i i3status-rs

# vi: ft=systemd
